# deploy/supervisor.conf
; -----------------------------------------------------------------------------
; supervisord 設定例(さくら VPS / Ubuntu 22.04 想定)
; /etc/supervisor/conf.d/unara.conf に配置する
; -----------------------------------------------------------------------------

[program:unara-py-worker]
command=/home/unara/worker/venv/bin/python /home/unara/worker/main.py
directory=/home/unara/worker
autostart=true
autorestart=true
user=unara
stdout_logfile=/var/log/unara/worker.log
stderr_logfile=/var/log/unara/worker.error.log
stopasgroup=true
killasgroup=true
environment=PYTHONUNBUFFERED="1"

; -----------------------------------------------------------------------------
; 注意: `unara:process-results` はワンショット型なので supervisor で常駐させない.
; autorestart=true にすると数百ms 単位の高速再起動ループになり Redis/MySQL に
; 過負荷が掛かる。常駐は cron の `schedule:run` (毎分) に一本化する.
; 即時消費が必要になったら、`unara:results-loop` を新規実装してから supervisor に
; 追加すること(while(true){ popResult; usleep(...); } の常駐型ハンドラ).
; -----------------------------------------------------------------------------
